# supervisor.conf
[program:pytech]
command=/home/olivx/www/pytech/start.sh
user=olivx
stdout_logfile=/home/olivx/www/pytech/logs/supervisor.log
redirect_stderr=true
autostart=true
autorestart=true
